# deploy/binance-trade-web.service
[Unit]
Description=binance-trade web dashboard (FastAPI + uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/binance-trade
# 绑定 127.0.0.1，仅经 nginx 反代对外；端口/凭据走 .env
Environment=WEB_HOST=127.0.0.1
Environment=WEB_PORT=8000
# -m web.server 需要项目根在 sys.path
Environment=PYTHONPATH=/opt/binance-trade
EnvironmentFile=-/opt/binance-trade/.env
ExecStart=/opt/binance-trade/.venv/bin/python -m web.server
Restart=on-failure
RestartSec=10
# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
# 允许写 SQLite/日志（项目目录）；ProtectSystem=full 会把 /usr /etc 设只读
ReadWritePaths=/opt/binance-trade
# 日志走 journald（journalctl -u binance-trade-web）。若要文件日志，
# 需 ReadWritePaths 包含日志目录后再用 StandardOutput=append:...

[Install]
WantedBy=multi-user.target
